# deploy/koc-engine-frontend.service
# KOC Engine — Frontend systemd service (Next.js)
# 安装: sudo cp koc-engine-frontend.service /etc/systemd/system/ && sudo systemctl daemon-reload && sudo systemctl enable --now koc-engine-frontend
# 管理: systemctl status/restart/stop koc-engine-frontend
# 日志: journalctl -u koc-engine-frontend -f

[Unit]
Description=KOC Engine Frontend (Next.js)
After=network.target
Wants=koc-engine.service

[Service]
Type=simple
User=koc-engine
Group=koc-engine
WorkingDirectory=/opt/koc-engine/frontend
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL=
Environment=PORT=3000
ExecStart=/usr/bin/node node_modules/.bin/next start -p 3000
Restart=always
RestartSec=3

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/koc-engine/frontend/.next/cache
ReadOnlyPaths=/opt/koc-engine/frontend

[Install]
WantedBy=multi-user.target
