# deploy/koc-engine.service
# KOC Engine — systemd service
# 安装: sudo cp koc-engine.service /etc/systemd/system/ && sudo systemctl daemon-reload && sudo systemctl enable --now koc-engine
# 管理: systemctl status/restart/stop koc-engine
# 日志: journalctl -u koc-engine -f

[Unit]
Description=KOC Engine API Server
After=network.target

[Service]
Type=simple
User=koc-engine
Group=koc-engine
WorkingDirectory=/opt/koc-engine/backend
Environment=PATH=/opt/koc-engine/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/koc-engine/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001 --workers 1
Restart=always
RestartSec=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/koc-engine/output
ReadOnlyPaths=/opt/koc-engine/backend

[Install]
WantedBy=multi-user.target
